# omserver.service
[Unit]
Description=OMServer

[Service]
ExecStart=/home/ati/runtestpythonV.sh

[Install]
WantedBy=multi-user.target

